# .setup/etc/systemd/system/greetd.service.d/quiet.conf
[Unit]
After=getty@tty2.service

[Service]
Type=idle
StandardInput=tty
TTYPath=/dev/tty2
TTYReset=yes
TTYVHangup=yes
